# mkosi.extra/usr/lib/systemd/system/mnt-nas.mount
[Unit]
Description=Mount NAS share from router SSD
After=network-online.target
Requires=network-online.target

[Mount]
What=//10.0.0.2/share
Where=/mnt/nas
Type=cifs
Options=guest,uid=1000,gid=1000,file_mode=0775,dir_mode=0775,nofail,_netdev

[Install]
WantedBy=multi-user.target
