# system_files/usr/lib/systemd/system/mnt-nas.automount
[Unit]
Description=Automount NAS share from router SSD

[Automount]
Where=/mnt/nas
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
